# raspberry/auto_start_file_1.service
[Unit]
Description=On-boot script launch servise

[Service]
ExecStart=/home/pi/auto_start_file.sh

[Install]
WantedBy=multi-user.target